# conf.d/server_supervisord.conf
[program: fastapi-server]
command=/root/  /start.sh
directory=/root/
autostart=true
autorestart=true
stderr_logfile=/root/  /out.log
stdout_logfile=/root/  /out.log
# Конфиг файл для супервизора